# roles/pinger/templates/pinger.service.j2
[Unit]
Description=Pinger service reporting connection latency
After=network-online.target

[Service]
User=pinger
ExecStart=/usr/local/bin/pinger

[Install]
WantedBy=multi-user.target